# etc/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
childlogdir=/var/log/

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[include]
files=/etc/supervisord/*.conf

[program:bot]
directory=/opt/bot/
command=/usr/bin/python deepanimebot/bot.py --debug
autorestart=true
startretries=10000
redirect_stderr=true
